# autossh/rootfs/etc/services.d/autossh/finish
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
bashio::log.warning "autossh service exited with code ${1:-?}; s6 will restart."
