# autossh/rootfs/etc/services.d/autossh/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
set -euo pipefail

# shellcheck source=/usr/lib/autossh/validate.sh
source /usr/lib/autossh/validate.sh
# shellcheck source=/usr/lib/autossh/keygen.sh
source /usr/lib/autossh/keygen.sh
# shellcheck source=/usr/lib/autossh/tunnel.sh
source /usr/lib/autossh/tunnel.sh

autossh::validate
autossh::keygen
autossh::tunnel
